# confs/supervisord.d/celery_beat.ini
[program:celery_beat]
command=/root/.virtualenvs/iyblog-_wyXyExX/bin/celery -A celery_worker:celery beat -l info
directory=/root/iyblog
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stderr_logfile=/var/log/celery/beat_err.log
stdout_logfile=/var/log/celery/beat_out.log
stopasgroup=true
killasgroup=true
