# confs/supervisord.d/celery_work.ini
[program:celery_work]
command=/root/.virtualenvs/iyblog-_wyXyExX/bin/celery -A celery_worker:celery worker -Q default,mail -c 10 -l info
directory=/root/iyblog
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stderr_logfile=/var/log/celery/worker_err.log
stdout_logfile=/var/log/celery/worker_out.log
stopasgroup=true
killasgroup=true
