# confs/supervisord.d/celery_work.ini
[program:celery_work]
command=/usr/bin/celery -A celery_worker:celery worker -Q default,hw_task -c 10 -l info
directory=/usr/local/PMMT      
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stderr_logfile=/var/log/celery/worker_err.log
stdout_logfile=/var/log/celery/worker_out.log
stopasgroup=true
killasgroup=true
